# scripts/soundboard-pipewire.service
[Unit]
Description=Sound Spring PipeWire virtual microphone routing
After=pipewire.service pipewire-pulse.service
Wants=pipewire.service pipewire-pulse.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%h/.local/share/soundspring/setup-pipewire.sh

[Install]
WantedBy=default.target
